# python-ml/crontab
59 23 * * * python /app/dns_blocker.py >> /var/log/cron.log 2>&1